# Microservices/microservices/services/view/run
#!/bin/bash

cd dosomedev-view

# Install project libraries.
npm install

# Run application in development mode.
npm run dev
